# elasticsearch/files/default/elastic.conf
[program:elastic]
directory=/mnt/elastic
command=/mnt/elastic/bin/elasticsearch
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/nodehook.err.log
stdout_logfile=/var/log/nodehook.out.log
user=elastic
